# deploy/xgrowth.service
# systemd unit for running without Docker.
# Install: see DEPLOY.md (the systemd path). Assumes a venv at
# /opt/autotwitter/.venv and .env at /opt/autotwitter/.env.

[Unit]
Description=X Growth Engine
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=xgrowth
WorkingDirectory=/opt/autotwitter
EnvironmentFile=/opt/autotwitter/.env
ExecStart=/opt/autotwitter/.venv/bin/python -m xgrowth.app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
